# config/myjob.conf
description     "berceuse"
start on startup
task
exec /root/Workshop/berceuse.sh
